# deploy/systemd/panaccess-wind.target
[Unit]
Description=PanAccess Wind — todas las instancias Daphne
Wants=panaccess-wind@8000.service

[Install]
WantedBy=multi-user.target
