# deploy/systemd/panaccess-wind@.service
[Unit]
Description=PanAccess Wind Daphne (puerto %i)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target postgresql.service redis-server.service
PartOf=panaccess-wind.target
# Ver panaccess-wind.service: el perfil de instancia única también usa el
# puerto 8000 -- Conflicts= evita que ambos perfiles corran a la vez en
# vez de depender de acordarse de deshabilitar uno a mano.
Conflicts=panaccess-wind.service

[Service]
Type=simple
User=wind
WorkingDirectory=/opt/panaccess-wind
EnvironmentFile=/opt/panaccess-wind/.env
ExecStart=/opt/panaccess-wind/env/bin/daphne -b 127.0.0.1 -p %i panaccess_wind_integration.asgi:application
Restart=always
RestartSec=5

[Install]
# Enlace directo a multi-user.target para arranque automático tras reboot
WantedBy=multi-user.target
