# deploy/systemd/panaccess-wind.service
[Unit]
Description=Servicio Web PanAccess Wind Integration (Daphne)
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=wind
WorkingDirectory=/opt/panaccess-wind
EnvironmentFile=/opt/panaccess-wind/.env
ExecStart=/opt/panaccess-wind/env/bin/daphne -b 127.0.0.1 -p 8000 panaccess_wind_integration.asgi:application
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
